# google-business-bot.service
[Unit]
Description=Google Business Bot
After=network.target
# Снимаем лимиты автоперезапуска, чтобы сервис самовосстанавливался без блокировок
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
# УКАЖИТЕ РЕАЛЬНЫЙ ПУТЬ К ПРОЕКТУ НА СЕРВЕРЕ
WorkingDirectory=/root/GoogleBusinessTelegramBot
# Экспортируем PATH, чтобы использовать интерпретатор из venv
Environment="PATH=/root/GoogleBusinessTelegramBot/new_venv/bin"
# УКАЖИТЕ РЕАЛЬНЫЙ ПУТЬ К .env
EnvironmentFile=/root/GoogleBusinessTelegramBot/.env
ExecStart=/root/GoogleBusinessTelegramBot/new_venv/bin/python /root/GoogleBusinessTelegramBot/bot.py

# Самовосстановление
Restart=always
RestartSec=10
TimeoutStartSec=180
TimeoutStopSec=30

# Логи в файлы (можно закомментировать и использовать journalctl)
StandardOutput=append:/root/GoogleBusinessTelegramBot/logs/bot.log
StandardError=append:/root/GoogleBusinessTelegramBot/logs/bot_error.log

# Создаем необходимые директории
ExecStartPre=/bin/mkdir -p /root/GoogleBusinessTelegramBot/logs
ExecStartPre=/bin/mkdir -p /root/GoogleBusinessTelegramBot/logs/context_logs

# Уведомление админу при падении, если systemd переведет сервис в failed
OnFailure=google-business-bot-notify@%n.service

[Install]
WantedBy=multi-user.target