# google-business-bot-update.service
[Unit]
Description=Google Business Bot KB Update (one-shot)
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/Users/test/Local_code/GoogleBusinessBot
Environment="PATH=/Users/test/Local_code/GoogleBusinessBot/new_venv/bin"
EnvironmentFile=/Users/test/Local_code/GoogleBusinessBot/.env
ExecStart=/Users/test/Local_code/GoogleBusinessBot/new_venv/bin/python /Users/test/Local_code/GoogleBusinessBot/update_kb.py
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal
OnFailure=google-business-bot-notify@%n.service

[Install]
WantedBy=multi-user.target


