# google-business-bot.service
[Unit]
Description=Google Business Bot
After=network.target
# Снимаем лимиты автоперезапуска, чтобы сервис самовосстанавливался без блокировок
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
# УКАЖИТЕ РЕАЛЬНЫЙ ПУТЬ К ПРОЕКТУ НА СЕРВЕРЕ
WorkingDirectory=/Users/test/Local_code/GoogleBusinessBot
# Экспортируем PATH, чтобы использовать интерпретатор из venv
Environment="PATH=/Users/test/Local_code/GoogleBusinessBot/new_venv/bin"
# УКАЖИТЕ РЕАЛЬНЫЙ ПУТЬ К .env
EnvironmentFile=/Users/test/Local_code/GoogleBusinessBot/.env
ExecStart=/Users/test/Local_code/GoogleBusinessBot/new_venv/bin/python /Users/test/Local_code/GoogleBusinessBot/bot.py

# Самовосстановление
Restart=always
RestartSec=10
TimeoutStartSec=180
TimeoutStopSec=30

# Логи в файлы (можно закомментировать и использовать journalctl)
StandardOutput=append:/Users/test/Local_code/GoogleBusinessBot/logs/bot.log
StandardError=append:/Users/test/Local_code/GoogleBusinessBot/logs/bot_error.log

# Создаем необходимые директории
ExecStartPre=/bin/mkdir -p /Users/test/Local_code/GoogleBusinessBot/logs
ExecStartPre=/bin/mkdir -p /Users/test/Local_code/GoogleBusinessBot/logs/context_logs

# Уведомление админу при падении, если systemd переведет сервис в failed
OnFailure=google-business-bot-notify@%n.service

[Install]
WantedBy=multi-user.target